# cron_ws/crontab
SHELL=/bin/bash

0 3 * * * bash ~/.config/cron/backup.sh
